# 99-scansnap-ix500.rules
# ScanSnap ix500 button polling service
# Start on connect, stop on disconnect
#
# Note: The remove action uses a workaround. If it doesn't work on your system,
# the service will handle disconnects gracefully by waiting and retrying.

# On connect - start the button polling service for logged-in user
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="04c5", ATTR{idProduct}=="132b", TAG+="systemd", ENV{SYSTEMD_USER_WANTS}="scan-button.service"

# On disconnect - the service handles this gracefully (waits and retries)
# Uncomment and adjust username if you want explicit stop:
# ACTION=="remove", SUBSYSTEM=="usb", ENV{ID_VENDOR_ID}=="04c5", ENV{ID_MODEL_ID}=="132b", RUN+="/usr/bin/systemctl --user --machine=YOUR_USERNAME@ stop scan-button.service"
